# etc/systemd/system/altrotate.timer
[Unit]
Description=Alternative Logrotate monitoring

[Timer]
OnBootSec=5m
OnUnitActiveSec=1d

[Install]
WantedBy=timers.target
