# etc/systemd/system/altrotate.service
[Unit]
Description=Alt Log Rotate

[Service]
Type=oneshot
ExecStart=/home/opal/SDKVM/scripts/blackbox/alt_logrotate.sh
